# supervisord.conf
[supervisord]
#nodaemon=true

[program:uwsgi]
command=/home/hicbrowser/conda/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --callable app

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

